# config/systemd/marketprism-websocket.service
[Unit]
Description=MarketPrism WebSocket守护服务
Documentation=https://github.com/MNS-Vic/marketprism
After=network.target
Wants=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/marketprism
Environment=NODE_ENV=production
Environment=PORT=8092
ExecStart=/usr/bin/python3 /home/ubuntu/marketprism/scripts/websocket-daemon-manager.py start
ExecStop=/usr/bin/python3 /home/ubuntu/marketprism/scripts/websocket-daemon-manager.py stop
ExecReload=/usr/bin/python3 /home/ubuntu/marketprism/scripts/websocket-daemon-manager.py restart

# 重启策略
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# 资源限制
MemoryLimit=512M
CPUQuota=50%

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/ubuntu/marketprism/logs

# 日志设置
StandardOutput=journal
StandardError=journal
SyslogIdentifier=marketprism-websocket

[Install]
WantedBy=multi-user.target
